# packaging/realestate-crm.service
[Unit]
Description=Real Estate CRM Web Server
Documentation=https://realestatecrm.com/docs
After=network.target network-online.target
Wants=network-online.target

[Service]
Type=simple
User=realestate-crm
Group=realestate-crm
WorkingDirectory=/opt/realestate-crm
Environment=PYTHONUNBUFFERED=1
Environment=API_HOST=0.0.0.0
Environment=API_PORT=6090
ExecStart=/opt/realestate-crm/.venv/bin/uvicorn backend.main:app --host 0.0.0.0 --port 6090
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=10
StandardOutput=append:/var/log/realestate-crm/server.log
StandardError=append:/var/log/realestate-crm/error.log
SyslogIdentifier=realestate-crm

# Security hardening
NoNewPrivileges=yes
ProtectSystem=full
ProtectHome=yes
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
ReadWritePaths=/var/lib/realestate-crm /var/log/realestate-crm /opt/realestate-crm

[Install]
WantedBy=multi-user.target
